# files/etc/uci-defaults/96-ath11k-mac80211-compat
#!/bin/sh
# Keep ath11k first-boot settings idempotent without patching runtime wireless helpers.

set -e

BOARD="$(cat /tmp/sysinfo/board_name 2>/dev/null || true)"

case "$BOARD" in
  jdcloud,re-ss-01|jdcloud_re-ss-01)
    ;;
  *)
    exit 0
    ;;
esac

for radio in radio0 radio1; do
  if uci -q get "wireless.${radio}" >/dev/null 2>&1; then
    if [ -z "$(uci -q get "wireless.${radio}.txantenna")" ]; then
      uci set "wireless.${radio}.txantenna=3"
    fi
    if [ -z "$(uci -q get "wireless.${radio}.rxantenna")" ]; then
      uci set "wireless.${radio}.rxantenna=3"
    fi
  fi
done

uci commit wireless >/dev/null 2>&1 || true

if [ -f /etc/init.d/rc-common-selfheal ]; then
  /bin/sh /rom/etc/rc.common /etc/init.d/rc-common-selfheal enable >/dev/null 2>&1 || true
fi

exit 0
