# files/etc/uci-defaults/96-ath11k-mac80211-compat
#!/bin/sh
# Work around ath11k drivers that reject iw distance/frag on first boot.

set -e

BOARD="$(cat /tmp/sysinfo/board_name 2>/dev/null || true)"
MAC80211_SH="/lib/netifd/wireless/mac80211.sh"

case "$BOARD" in
  jdcloud,re-ss-01|jdcloud_re-ss-01)
    ;;
  *)
    rm -f "$0"
    exit 0
    ;;
esac

if [ -f "$MAC80211_SH" ]; then
  sed -i '/set distance/d' "$MAC80211_SH"
  sed -i '/set frag/d' "$MAC80211_SH"
fi

for radio in radio0 radio1; do
  if uci -q get "wireless.${radio}" >/dev/null 2>&1; then
    if [ -z "$(uci -q get "wireless.${radio}.txantenna")" ]; then
      uci set "wireless.${radio}.txantenna=3"
    fi
    if [ -z "$(uci -q get "wireless.${radio}.rxantenna")" ]; then
      uci set "wireless.${radio}.rxantenna=3"
    fi
  fi
done

uci commit wireless >/dev/null 2>&1 || true

if [ -f /etc/init.d/rc-common-selfheal ]; then
  /bin/sh /rom/etc/rc.common /etc/init.d/rc-common-selfheal enable >/dev/null 2>&1 || true
fi

rm -f "$0"
exit 0
