# files/etc/uci-defaults/98-home-partition
#!/bin/sh
# /home 分区自动挂载脚本
# 优先复用 GPT 中名为 "primary" 的大数据分区；
# 不再假定 p22/p21 固定可用，避免误用系统分区。

set -e

DISK='/dev/mmcblk0'
PART=''
LEGACY_PART='/dev/mmcblk0p21'
TARGET='/home'
PXE_DEFAULTS='/usr/share/mgrserver-defaults/res/pxe'
PXE_TARGET='/home/pxe'
CREATED_PART='0'
HOME_ON_OVERLAY='0'
MIN_HOME_PART_MIB='4096'
RC_COMMON='/rom/etc/rc.common'

run_init_script() {
    local name="$1"
    shift
    local script="/etc/init.d/$name"

    [ -f "$script" ] || return 127

    if [ -f "$RC_COMMON" ]; then
        /bin/sh "$RC_COMMON" "$script" "$@"
        return $?
    fi

    "$script" "$@"
}

run_with_timeout() {
    local seconds="$1"
    shift

    if command -v timeout >/dev/null 2>&1; then
        timeout "$seconds" "$@"
        return $?
    fi

    "$@"
}

part_size_mib() {
    awk -v name="${1##*/}" '$4 == name {print int($3 / 1024)}' /proc/partitions 2>/dev/null | tail -n1
}

part_name() {
    awk -F= '$1 == "PARTNAME" {print $2}' "/sys/class/block/${1##*/}/uevent" 2>/dev/null | tail -n1
}

is_reserved_part_name() {
    case "$1" in
        ''|rootfs|rootfs_1|rootfs_data|primary)
            return 1
            ;;
        0:*|HLOS|HLOS_1|WIFIFW|WIFIFW_1|BOOTCONFIG|BOOTCONFIG1|QSEE|QSEE_1|DEVCFG|DEVCFG_1|RPM|RPM_1|CDT|CDT_1|APPSBLENV|APPSBL|APPSBL_1|ART|ETHPHYFW|plugin|log|swap)
            return 0
            ;;
    esac

    return 1
}

find_primary_part() {
    for SYS_PART in /sys/class/block/mmcblk0/mmcblk0p*; do
        [ -e "$SYS_PART" ] || continue
        CANDIDATE="/dev/$(basename "$SYS_PART")"
        [ -b "$CANDIDATE" ] || continue

        PART_SIZE_MIB=$(part_size_mib "$CANDIDATE")
        [ -n "$PART_SIZE_MIB" ] || continue
        [ "$PART_SIZE_MIB" -gt "$MIN_HOME_PART_MIB" ] || continue

        PART_NAME=$(part_name "$CANDIDATE")
        if [ "$PART_NAME" = 'primary' ]; then
            echo "$CANDIDATE"
            return 0
        fi
    done

    return 1
}

find_largest_safe_part() {
    BEST_PART=''
    BEST_NUM='0'

    for SYS_PART in /sys/class/block/mmcblk0/mmcblk0p*; do
        [ -e "$SYS_PART" ] || continue
        PART_BASENAME=$(basename "$SYS_PART")
        PART_NUM=${PART_BASENAME#mmcblk0p}
        case "$PART_NUM" in
            ''|*[!0-9]*)
                continue
                ;;
        esac

        CANDIDATE="/dev/$PART_BASENAME"
        [ -b "$CANDIDATE" ] || continue

        PART_SIZE_MIB=$(part_size_mib "$CANDIDATE")
        [ -n "$PART_SIZE_MIB" ] || continue
        [ "$PART_SIZE_MIB" -gt "$MIN_HOME_PART_MIB" ] || continue

        PART_NAME=$(part_name "$CANDIDATE")
        if is_reserved_part_name "$PART_NAME"; then
            continue
        fi

        if [ "$PART_NUM" -gt "$BEST_NUM" ]; then
            BEST_NUM="$PART_NUM"
            BEST_PART="$CANDIDATE"
        fi
    done

    [ -n "$BEST_PART" ] || return 1
    echo "$BEST_PART"
    return 0
}

find_home_part() {
    FOUND_PART=$(find_primary_part || true)
    if [ -n "$FOUND_PART" ]; then
        echo "$FOUND_PART"
        return 0
    fi

    if [ -b "$LEGACY_PART" ]; then
        LEGACY_SIZE_MIB=$(part_size_mib "$LEGACY_PART")
        if [ -n "$LEGACY_SIZE_MIB" ] && [ "$LEGACY_SIZE_MIB" -gt "$MIN_HOME_PART_MIB" ]; then
            echo "$LEGACY_PART"
            return 0
        fi
    fi

    find_largest_safe_part
}

sync_pxe_defaults() {
    [ -d "$PXE_DEFAULTS" ] || return 0

    mkdir -p "$PXE_TARGET"

    if [ -f "$PXE_TARGET/config.ini" ] && [ -f "$PXE_TARGET/up_pxe_res.sh" ] && [ -d "$PXE_TARGET/tftpboot" ]; then
        echo "PXE defaults already present: $PXE_TARGET"
        return 0
    fi

    cp -a "$PXE_DEFAULTS"/. "$PXE_TARGET"/ 2>/dev/null \
        || { echo "Failed to seed PXE defaults from $PXE_DEFAULTS"; return 1; }

    echo "PXE defaults restored: $PXE_DEFAULTS -> $PXE_TARGET"
    return 0
}

apply_pxe_setup() {
    local script="$PXE_TARGET/set_pxe.sh"

    [ -f "$script" ] || { echo "PXE setup skipped: $script missing"; return 0; }
    chmod +x "$script" 2>/dev/null || true

    OFFLINE_STRICT=true QUIET=1 /bin/sh "$script" \
        || { echo "PXE setup failed: $script"; return 1; }

    echo "PXE setup applied: $script"
    return 0
}

mounted_source_for_target() {
    awk -v target="$1" '$2 == target {print $1; exit}' /proc/mounts 2>/dev/null
}

unmount_duplicate_part_mounts() {
    local device="$1"
    local target="$2"

    awk -v device="$device" -v target="$target" '$1 == device && $2 != target {print $2}' /proc/mounts 2>/dev/null |
        while read -r mount_point; do
            [ -n "$mount_point" ] || continue
            umount "$mount_point" 2>/dev/null || true
        done
}

# 仅在 eMMC 设备存在时执行
if [ ! -b "$DISK" ]; then
    echo "Disk device not found: $DISK, falling back to overlay /home"
    mkdir -p "$TARGET"
    HOME_ON_OVERLAY='1'
fi

# 优先使用 GPT 中名为 primary 的大分区，其次才是 legacy 布局
if [ "$HOME_ON_OVERLAY" != '1' ]; then
    PART=$(find_home_part || true)
fi

# 如果分区不存在，创建它
if [ "$HOME_ON_OVERLAY" != '1' ] && [ -z "$PART" ]; then
    echo 'No safe /home partition detected, creating a new primary data partition'
    # 仅使用 parted，避免 fdisk 分区号自动选择的歧义
    if command -v parted >/dev/null 2>&1; then
        PART_START_MIB=$(parted -ms "$DISK" unit MiB print 2>/dev/null \
            | awk -F: '$1 ~ /^[0-9]+$/ {gsub("MiB", "", $3); end=$3} END {if (end == "") print ""; else printf "%d", int(end) + 1}')
        [ -n "$PART_START_MIB" ] || PART_START_MIB='1'
        run_with_timeout 30 parted -s "$DISK" mkpart primary ext4 "${PART_START_MIB}MiB" 100% \
            || { echo "parted failed, falling back to overlay /home"; mkdir -p "$TARGET"; HOME_ON_OVERLAY='1'; }
    else
        echo 'parted not available, falling back to overlay /home.'
        mkdir -p "$TARGET"
        HOME_ON_OVERLAY='1'
    fi
    if [ "$HOME_ON_OVERLAY" != '1' ]; then
        command -v partprobe >/dev/null 2>&1 && partprobe "$DISK" || true
        sleep 2
        PART=$(find_home_part || true)
        CREATED_PART='1'
    fi
fi

# 等待设备节点出现（parted 创建后可能需要时间刷新）
if [ "$HOME_ON_OVERLAY" != '1' ] && [ -n "$PART" ] && [ ! -b "$PART" ]; then
    mdev -s >/dev/null 2>&1 || true
    sleep 2
    # 再次检查是否出现（解决 parted 后立即检测失败的问题）
    [ -b "$PART" ] || {
        echo "Partition node not appeared after creation: $PART, falling back to overlay /home"
        mkdir -p "$TARGET"
        HOME_ON_OVERLAY='1'
    }
fi

if [ "$HOME_ON_OVERLAY" != '1' ] && [ -z "$PART" ]; then
    echo 'Unable to resolve a safe /home partition after partitioning, falling back to overlay /home'
    mkdir -p "$TARGET"
    HOME_ON_OVERLAY='1'
fi

# 格式化（新分区或无文件系统时）
PART_FSTYPE='overlay'
if [ "$HOME_ON_OVERLAY" != '1' ]; then
    command -v blkid >/dev/null 2>&1 || {
        echo 'blkid not found, falling back to overlay /home'
        mkdir -p "$TARGET"
        HOME_ON_OVERLAY='1'
    }
fi
if [ "$HOME_ON_OVERLAY" != '1' ]; then
    PART_FSTYPE=$(blkid -s TYPE -o value "$PART" 2>/dev/null || true)
    if [ "$CREATED_PART" = '1' ] || [ -z "$PART_FSTYPE" ]; then
        command -v mkfs.ext4 >/dev/null 2>&1 || {
            echo 'mkfs.ext4 not found, falling back to overlay /home'
            mkdir -p "$TARGET"
            HOME_ON_OVERLAY='1'
        }
    fi
fi
if [ "$HOME_ON_OVERLAY" != '1' ] && { [ "$CREATED_PART" = '1' ] || [ -z "$PART_FSTYPE" ]; }; then
    run_with_timeout 60 mkfs.ext4 -F "$PART" >/dev/null \
        || { echo "mkfs.ext4 failed on $PART, falling back to overlay /home"; mkdir -p "$TARGET"; HOME_ON_OVERLAY='1'; }
    PART_FSTYPE='ext4'
fi

# 挂载
mkdir -p "$TARGET"
SKIP_MOUNT='0'
if [ "$HOME_ON_OVERLAY" != '1' ]; then
    CURR_SRC=$(mounted_source_for_target "$TARGET")
fi
if [ "$HOME_ON_OVERLAY" != '1' ] && [ -n "${CURR_SRC:-}" ]; then
    if [ "$CURR_SRC" = "$PART" ]; then
        SKIP_MOUNT='1'
    else
        umount "$TARGET" || true
    fi
fi
if [ "$HOME_ON_OVERLAY" != '1' ] && [ "$SKIP_MOUNT" != '1' ]; then
    run_with_timeout 10 mount "$PART" "$TARGET" 2>/dev/null \
        || run_with_timeout 10 mount -t "$PART_FSTYPE" "$PART" "$TARGET" \
        || { echo "mount failed, falling back to overlay /home"; mkdir -p "$TARGET"; HOME_ON_OVERLAY='1'; }
fi
if [ "$HOME_ON_OVERLAY" != '1' ]; then
    unmount_duplicate_part_mounts "$PART" "$TARGET"
fi

# 写入 fstab 持久化（使用双引号使变量正确展开）
if [ "$HOME_ON_OVERLAY" != '1' ]; then
    [ -f /etc/config/fstab ] || touch /etc/config/fstab
    uci -q delete fstab.mgrhome || true
    uci set fstab.mgrhome='mount'
    uci set fstab.mgrhome.target="$TARGET"
    uci set fstab.mgrhome.device="$PART"
    uci set fstab.mgrhome.fstype="$PART_FSTYPE"
    uci set fstab.mgrhome.options='rw,sync'
    uci set fstab.mgrhome.enabled='1'
    uci set fstab.mgrhome.enabled_fsck='0'
    uci commit fstab
    run_init_script fstab enable >/dev/null 2>&1 || true
fi

# 将 MgrServer 默认工具文件（启动器等）从固件只读区同步到 /home/tools
sync_tools_defaults() {
    local src="/usr/share/mgrserver-defaults/res"
    local dst="/home/tools/launcher"

    [ -d "$src" ] || { echo "Tools source directory missing: $src"; return 0; }

    mkdir -p "$dst"

    # launcher
    local launcher_src="$src/yushu-tool-launcher.exe"
    local launcher_dst="$dst/yushu-tool-launcher.exe"

    if [ -f "$launcher_src" ] && { [ ! -f "$launcher_dst" ] || [ ! -s "$launcher_dst" ]; }; then
        cp "$launcher_src" "$launcher_dst" 2>/dev/null \
            || { echo "Failed to copy launcher from $launcher_src"; return 1; }
        echo "Launcher seeded: $launcher_src -> $launcher_dst"
    elif [ -f "$launcher_dst" ]; then
        echo "Launcher already present: $launcher_dst"
    fi

    return 0
}

# 仅在 /home 分区正常挂载时同步默认文件（overlay 回退路径会被丢弃，不应写入）
if [ "$HOME_ON_OVERLAY" != '1' ]; then
    sync_pxe_defaults || true
    apply_pxe_setup || true
    sync_tools_defaults || true
fi

if [ "$HOME_ON_OVERLAY" = '1' ]; then
    echo "/home fallback ready on overlay: $TARGET (will retry partition setup on next boot)"
    exit 1
fi

echo "/home partition ready: $PART -> $TARGET ($PART_FSTYPE)"
exit 0
