# files/etc/uci-defaults/99-mgrserver-ports
#!/bin/sh
# MgrServer 端口配置 + LuCI 端口迁移
# 对齐 deploy-openwrt.sh Step 5.1-5.3 逻辑

set -e

. /lib/functions.sh

MGR_PORT="80"
LUCI_PORT="3500"
SCRIPT_OK="0"
LOG_TAG="mgrserver-ports"

cleanup() {
  [ "$SCRIPT_OK" = "1" ] && rm -f "$0"
}

trap cleanup EXIT

log_msg() {
  logger -t "$LOG_TAG" "$1" 2>/dev/null || true
}

rewrite_nginx_listen_value() {
  local entry="$1"
  local first_token=""
  local rest=""
  local rewritten_first=""

  [ -n "$entry" ] || return 1

  first_token="${entry%% *}"
  if [ "$first_token" = "$entry" ]; then
    rest=""
  else
    rest="${entry#"$first_token"}"
  fi

  case "$first_token" in
    80)
      rewritten_first="$LUCI_PORT"
      ;;
    \[*\]:80|*:80|[A-Za-z0-9._-]*:80)
      rewritten_first="${first_token%:80}:$LUCI_PORT"
      ;;
    *)
      rewritten_first="$first_token"
      ;;
  esac

  if [ "$rewritten_first" != "$first_token" ]; then
    printf '%s%s' "$rewritten_first" "$rest"
    return 0
  fi

  return 1
}

ensure_mgrserver_started() {
  local attempt

  /etc/init.d/mgrserver enable >/dev/null 2>&1 || true
  /etc/init.d/mgrserver restart >/dev/null 2>&1 || /etc/init.d/mgrserver start >/dev/null 2>&1

  for attempt in 1 2 3 4 5; do
    if /etc/init.d/mgrserver health_check >/dev/null 2>&1; then
      log_msg "MgrServer is healthy on port ${MGR_PORT}"
      return 0
    fi
    sleep 2
  done

  log_msg "MgrServer failed health check after port migration"
  return 1
}

# ============================================================
# MgrServer 基础配置
# ============================================================

# 确保 /etc/mgrserver 目录存在 (deploy 脚本通过 SSH mkdir 创建)
mkdir -p /etc/mgrserver
chmod 700 /etc/mgrserver

if [ ! -f /etc/config/mgrserver ]; then
  touch /etc/config/mgrserver
fi
uci -q get mgrserver.main >/dev/null 2>&1 || uci set mgrserver.main=mgrserver
uci set mgrserver.main.port="$MGR_PORT"
uci commit mgrserver

# ============================================================
# 端口冲突处理: MgrServer→:80, LuCI→:3500
# 与 deploy-openwrt.sh Step 5.1-5.3 逻辑对齐
# ============================================================

NGINX_ENABLED=0
if [ -x /etc/init.d/nginx ] && /etc/init.d/nginx enabled >/dev/null 2>&1; then
  NGINX_ENABLED=1
fi

UHTTPD_SECTIONS="$(uci show uhttpd 2>/dev/null | awk -F'[.=]' '/=uhttpd$/ {print $2}')"
if [ -z "$UHTTPD_SECTIONS" ]; then
  UHTTPD_SECTION="$(uci add uhttpd uhttpd)"
  UHTTPD_SECTIONS="$UHTTPD_SECTION"
else
  if echo "$UHTTPD_SECTIONS" | grep -qx main; then
    UHTTPD_SECTION="main"
  else
    UHTTPD_SECTION="$(echo "$UHTTPD_SECTIONS" | head -n1)"
  fi
fi

for sec in $UHTTPD_SECTIONS; do
  uci -q del_list uhttpd."$sec".listen_http="0.0.0.0:80"
  uci -q del_list uhttpd."$sec".listen_http="[::]:80"
  uci -q del_list uhttpd."$sec".listen_https="0.0.0.0:443"
  uci -q del_list uhttpd."$sec".listen_https="[::]:443"
  if [ "$NGINX_ENABLED" = "1" ]; then
    uci -q del_list uhttpd."$sec".listen_http="0.0.0.0:$LUCI_PORT"
    uci -q del_list uhttpd."$sec".listen_http="[::]:$LUCI_PORT"
    uci -q del_list uhttpd."$sec".listen_https="0.0.0.0:$LUCI_PORT"
    uci -q del_list uhttpd."$sec".listen_https="[::]:$LUCI_PORT"
  fi
done

if [ "$NGINX_ENABLED" != "1" ]; then
  uci set uhttpd."$UHTTPD_SECTION".enabled="1"
  uci -q delete uhttpd."$UHTTPD_SECTION".listen_http
  uci -q delete uhttpd."$UHTTPD_SECTION".listen_https
  uci add_list uhttpd."$UHTTPD_SECTION".listen_http="0.0.0.0:$LUCI_PORT"
  uci add_list uhttpd."$UHTTPD_SECTION".listen_http="[::]:$LUCI_PORT"
  uci add_list uhttpd."$UHTTPD_SECTION".listen_https="0.0.0.0:$LUCI_PORT"
  uci add_list uhttpd."$UHTTPD_SECTION".listen_https="[::]:$LUCI_PORT"
fi

uci commit uhttpd
/etc/init.d/uhttpd restart >/dev/null 2>&1 || true

if [ "$NGINX_ENABLED" = "1" ]; then
  rewrite_nginx_listen_entry() {
    local entry="$1"
    local new_entry=""

    [ -n "$entry" ] || return 0

    if new_entry="$(rewrite_nginx_listen_value "$entry")"; then
      NGINX_SECTION_CHANGED=1
    else
      new_entry="$entry"
    fi

    NGINX_NEW_LISTENS="${NGINX_NEW_LISTENS}${new_entry}
"
  }

  rewrite_nginx_section() {
    local section="$1"
    local new_entry=""

    NGINX_NEW_LISTENS=""
    NGINX_SECTION_CHANGED=0
    config_list_foreach "$section" listen rewrite_nginx_listen_entry

    if [ "$NGINX_SECTION_CHANGED" = "1" ]; then
      NGINX_CHANGED=1
      uci -q delete nginx."$section".listen
      while IFS= read -r new_entry; do
        [ -n "$new_entry" ] || continue
        uci add_list nginx."$section".listen="$new_entry"
      done <<EOF_LISTEN_DST
$NGINX_NEW_LISTENS
EOF_LISTEN_DST
    fi
  }

  config_load nginx
  NGINX_SECTIONS="$(uci show nginx 2>/dev/null | awk -F'[.=]' '/=server$/ {print $2}')"
  NGINX_CHANGED=0
  for sec in $NGINX_SECTIONS; do
    rewrite_nginx_section "$sec"
  done

  if [ "$NGINX_CHANGED" = "1" ]; then
    uci commit nginx
    /etc/init.d/nginx restart >/dev/null 2>&1 || /etc/init.d/nginx start >/dev/null 2>&1
  fi
fi

# ============================================================
# 启用 MgrServer 服务
# ============================================================
ensure_mgrserver_started

SCRIPT_OK="1"

exit 0
