# files/etc/uci-defaults/99-service-watchdog-cron
#!/bin/sh
# Setup cron job for service watchdog
# Runs every 2 minutes to check service health

RC_COMMON="/rom/etc/rc.common"

run_init_script() {
  local name="$1"
  shift
  local script="/etc/init.d/$name"

  [ -f "$script" ] || return 127

  if [ -f "$RC_COMMON" ]; then
    /bin/sh "$RC_COMMON" "$script" "$@"
    return $?
  fi

  "$script" "$@"
}

# Add cron job for service watchdog
if command -v crontab >/dev/null 2>&1; then
  (crontab -l 2>/dev/null | grep -v "service-watchdog"; echo "*/2 * * * * /usr/bin/service-watchdog") | crontab -
else
  mkdir -p /etc/crontabs
  touch /etc/crontabs/root
  grep -v "service-watchdog" /etc/crontabs/root > /tmp/mgrserver-cron.$$ 2>/dev/null || true
  echo "*/2 * * * * /usr/bin/service-watchdog" >> /tmp/mgrserver-cron.$$
  mv /tmp/mgrserver-cron.$$ /etc/crontabs/root
fi

# Restart cron
if [ -x /etc/init.d/cron ]; then
  run_init_script cron restart >/dev/null 2>&1 || run_init_script cron reload >/dev/null 2>&1 || true
elif command -v crond >/dev/null 2>&1; then
  killall -HUP crond >/dev/null 2>&1 || crond -c /etc/crontabs >/dev/null 2>&1 || true
fi

exit 0
