# files/etc/uci-defaults/99-service-watchdog-cron
#!/bin/sh
# Setup cron job for service watchdog
# Runs every 5 minutes to check service health

# Add cron job for service watchdog
if command -v crontab >/dev/null 2>&1; then
  (crontab -l 2>/dev/null | grep -v "service-watchdog"; echo "*/5 * * * * /usr/bin/service-watchdog") | crontab -
else
  mkdir -p /etc/crontabs
  touch /etc/crontabs/root
  grep -v "service-watchdog" /etc/crontabs/root > /tmp/mgrserver-cron.$$ 2>/dev/null || true
  echo "*/5 * * * * /usr/bin/service-watchdog" >> /tmp/mgrserver-cron.$$
  mv /tmp/mgrserver-cron.$$ /etc/crontabs/root
fi

# Restart cron
if [ -x /etc/init.d/cron ]; then
  /etc/init.d/cron restart >/dev/null 2>&1 || /etc/init.d/cron reload >/dev/null 2>&1 || true
elif command -v crond >/dev/null 2>&1; then
  killall -HUP crond >/dev/null 2>&1 || crond -c /etc/crontabs >/dev/null 2>&1 || true
fi

exit 0
