# files/etc/hotplug.d/iface/30-auto-lan-subnet
#!/bin/sh
# Auto-adjust LAN subnet when it conflicts with WAN subnet.
# Triggered by hotplug on wan ifup events.
#
# When the upstream router's LAN and this router's LAN are on the
# same subnet (e.g., both 192.168.1.0/24), IPv4 routing breaks
# because the kernel can't disambiguate routes.  This script detects
# the conflict and moves the LAN to a non-conflicting /24 subnet.

[ "$ACTION" = "ifup" ] || exit 0
[ "$INTERFACE" = "wan" ] || exit 0

log_msg() {
  logger -t "auto-lan-subnet" "$1" 2>/dev/null || true
}

# --- 1. Get WAN IPv4 CIDR (e.g. "192.168.1.7/24") ---
wan_cidr=$(ip -4 -o addr show dev wan 2>/dev/null | awk '{print $4}')
[ -n "$wan_cidr" ] || exit 0

# --- 2. Get current LAN config from UCI ---
lan_ip_cidr=$(uci -q get network.lan.ipaddr 2>/dev/null || echo "")
[ -n "$lan_ip_cidr" ] || exit 0

# --- 3. Extract /24 prefix for comparison ---
# Handle both "x.x.x.x" and "x.x.x.x/yy" formats
wan_ip="${wan_cidr%%/*}"
wan_prefix=$(echo "$wan_ip" | cut -d. -f1-3)

lan_ip="${lan_ip_cidr%%/*}"
lan_prefix=$(echo "$lan_ip" | cut -d. -f1-3)

# Different /24 blocks → no conflict
[ "$wan_prefix" = "$lan_prefix" ] || exit 0

# --- 4. Conflict detected — find a non-conflicting /24 ---
new_third=2
while [ "$new_third" -le 200 ]; do
  candidate="192.168.${new_third}"
  [ "$candidate" = "$wan_prefix" ] || break
  new_third=$((new_third + 1))
done

[ "$new_third" -le 200 ] || {
  log_msg "no available /24 for LAN (all 192.168.x.0/24 blocked by WAN subnet $wan_prefix)"
  exit 1
}

new_lan_ip="192.168.${new_third}.1/24"

# --- 5. Apply ---
log_msg "subnet conflict: WAN=$wan_cidr LAN=$lan_ip_cidr → moving LAN to $new_lan_ip"

uci set network.lan.ipaddr="$new_lan_ip"
uci commit network

# Reload is gentler than restart — applies config changes without
# tearing down unrelated interfaces.
/etc/init.d/network reload 2>/dev/null || true

log_msg "LAN subnet adjusted to $new_lan_ip"

# --- 6. Broadcast gratuitous ARP so connected clients pick up the
#        new gateway quickly (DHCP renewal can take minutes).
sleep 1
arping -c 3 -I br-lan -U "${new_lan_ip%%/*}" 2>/dev/null || true

exit 0
